# pagerduty.service
[Unit]
Description=PagerDuty Incident Dashboard
After=network.target

[Service]
Type=simple
User=pagerduty
Group=pagerduty
WorkingDirectory=/opt/pagerduty
ExecStart=/bin/bash -c '. "$HOME/.nvm/nvm.sh" 2>/dev/null; exec npm --prefix /opt/pagerduty start'
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pagerduty
EnvironmentFile=-/opt/pagerduty/.env.local

[Install]
WantedBy=multi-user.target
