# pagerduty.service
[Unit]
Description=PagerDuty Incident Dashboard
After=network.target

[Service]
Type=simple
User=pagerduty
Group=pagerduty
WorkingDirectory=/opt/pagerduty
Environment=HOME=/home/pagerduty
ExecStart=/bin/bash -c 'if [ -f "$HOME/.nvm/nvm.sh" ]; then source "$HOME/.nvm/nvm.sh"; fi; command -v npm >/dev/null || { echo "npm not found for pagerduty user"; exit 1; }; exec npm start'
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pagerduty
EnvironmentFile=-/opt/pagerduty/.env.local

[Install]
WantedBy=multi-user.target
